# systemd/ukwinika-realtime-backup.service
[Unit]
Description=UKwinika Real-Time File Monitoring Service (v2.3 - inotify)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/enhanced_automated_backups.sh real-time
Restart=always
RestartSec=10
Nice=19
User=root
