# systemd/ukwinika-backup.service
[Unit]
Description=UKwinika Daily Idempotent Backup (v3.0)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Config and secrets paths are injected as environment variables
Environment="UKW_CONFIG=/etc/ukwinika-backup.conf"
Environment="UKW_SECRETS=/etc/ukwinika-backup.secrets"
ExecStart=/usr/local/bin/enhanced_automated_backups.sh backup
User=root
Nice=19
IOSchedulingClass=idle
PrivateTmp=true
