# systemd/ukwinika-backup.service
[Unit]
Description=UKwinika Daily Backup Service
After=network.target remote-fs.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/enhanced_automated_backups.sh backup incremental borg
EnvironmentFile=-/etc/ukwinika-backup.conf
User=root

# Nice level and I/O priority for production use
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=true
