# systemd/ukwinika-realtime-backup.service
# =============================================================================
# UKwinika Real-Time File Monitoring Service - v2.4
# Uses inotify for instant backup on file changes
# =============================================================================
[Unit]
Description=UKwinika Real-Time Backup Service (v2.4 - inotify monitoring)
After=network.target remote-fs.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/enhanced_automated_backups.sh real-time
Restart=always
RestartSec=5
EnvironmentFile=-/etc/ukwinika-backup.conf
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
